# sdk_container/src/third_party/coreos-overlay/coreos/sysext/devext/usr/lib/systemd/system/init-devext.service
[Unit]
Description=Flatcar developer system extension

[Service]
Type=oneshot
ExecStart=/usr/share/flatcar/devext/devext-init.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
