# sdk_container/src/third_party/portage-stable/sys-fs/mdadm/files/mdadm.timer
[Unit]
Description=Weekly check for MD array's redundancy information.

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=weekly
Persistent=true
