# sdk_container/src/third_party/coreos-overlay/coreos-base/oem-azure/files/units/waagent.service
[Unit]
Description=Microsoft Azure Agent
Wants=network-online.target sshd-keygen.service
After=network-online.target sshd-keygen.service

[Service]
Type=simple
Restart=always
RestartSec=5s
Environment=PATH=/oem/python/bin:/oem/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONUNBUFFERED=x
ExecStart=/oem/bin/waagent -daemon

[Install]
WantedBy=multi-user.target
