# sdk_container/src/third_party/coreos-overlay/coreos/sysext/devext/usr/lib/systemd/system/systemd-sysext.service.d/10-devext.conf
[Service]
ExecStart=
ExecStart=systemd-sysext --mutable=auto refresh
ExecReload=
ExecReload=systemd-sysext --mutable=auto refresh

