# sdk_container/src/third_party/portage-stable/app-admin/rsyslog/files/3-stable/rsyslog.init
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/rsyslog/files/3-stable/rsyslog.init,v 1.2 2008/08/14 06:47:38 dev-zero Exp $

opts="reload"

depend() {
	provide logger
}

start() {
	ebegin "Starting rsyslogd"
	start-stop-daemon \
		--start --quiet --exec /usr/sbin/rsyslogd \
		-- ${RSYSLOG_OPTS} -i "${PIDFILE}" -f "${CONFIGFILE}"
	eend $?
}

stop() {
	ebegin "Stopping rsyslogd"
	start-stop-daemon \
		--stop --quiet \
		--pidfile "${PIDFILE}"
	eend $?
}

reload() {
	if [ ! -f "${PIDFILE}" ]; then
		eerror "rsyslogd not running"
		return 1
	fi
	
	ebegin "Reloading configuration and re-opening log files"
	start-stop-daemon --stop --oknodo --signal HUP \
		--pidfile "${PIDFILE}"
	eend $?
}
