# sdk_container/src/third_party/coreos-overlay/app-emulation/hv-daemons/files/hv_fcopy_daemon.service
[Unit]
Description=Hyper-V FCOPY daemon
ConditionPathExists=/dev/vmbus/hv_fcopy

[Service]
ExecStart=/usr/bin/hv_fcopy_daemon --no-daemon

[Install]
WantedBy=multi-user.target
