# sdk_container/src/third_party/coreos-overlay/app-emulation/hv-daemons/files/hv_kvp_daemon.service
[Unit]
Description=Hyper-V KVP daemon
ConditionPathExists=/dev/vmbus/hv_kvp

[Service]
ExecStart=/usr/bin/hv_kvp_daemon --no-daemon

[Install]
WantedBy=multi-user.target
