# sdk_container/src/third_party/coreos-overlay/app-emulation/hv-daemons/files/hv_kvp_daemon.service
[Unit]
Description=Hyper-V KVP daemon
BindsTo=sys-devices-virtual-misc-vmbus\x21hv_kvp.device
ConditionPathExists=/dev/vmbus/hv_kvp

[Service]
Type=exec
ExecStart=/usr/bin/hv_kvp_daemon --no-daemon

[Install]
WantedBy=multi-user.target
