# sdk_container/src/third_party/coreos-overlay/app-emulation/hv-daemons/files/hv_vss_daemon.service
[Unit]
Description=Hyper-V VSS daemon
ConditionPathExists=/dev/vmbus/hv_vss

[Service]
ExecStart=/usr/bin/hv_vss_daemon --no-daemon

[Install]
WantedBy=multi-user.target
