# repos/flatcar-overlay/coreos-base/oem-gce/files/units/oem-gce.service
[Unit]
Description=GCE Linux Agent
After=local-fs.target network-online.target

[Service]
Type=notify
NotifyAccess=all
Restart=always
RestartSec=5

# There is a custom main process that kills all of the contained services.
KillMode=process
KillSignal=SIGTERM

ExecStart=/usr/bin/init.sh

[Install]
WantedBy=multi-user.target
