# sdk_container/src/third_party/coreos-overlay/sys-fs/mdadm/files/mdadm.service
[Unit]
Description=Initiates a check run of an MD array's redundancy information.

[Service]
Type=oneshot
ExecStart=/usr/sbin/checkarray --cron --all --idle --quiet
