# sdk_container/src/third_party/coreos-overlay/app-admin/logrotate/files/logrotate.service
[Unit]
Description=Rotate and Compress System Logs

[Service]
ExecStart=/usr/bin/logrotate /usr/share/logrotate/logrotate.conf
