# sdk_container/src/third_party/coreos-overlay/coreos-base/oem-stackit/files/ptp-kvm-load.service
[Unit]
Description=Load ptp_kvm kernel module
After=local-fs.target network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/modprobe ptp_kvm

[Install]
WantedBy=multi-user.target
