# sdk_container/src/third_party/coreos-overlay/x11-drivers/nvidia-drivers/files/units/nvidia.service
[Unit]
Description=NVIDIA Configure Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=no
Environment=PATH=/opt/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/usr/share/oem/bin/setup-nvidia

[Install]
WantedBy=multi-user.target
