# sdk_container/src/third_party/coreos-overlay/coreos-base/oem-gce/files/units/oem-gce.service
[Unit]
Description=GCE Linux Agent
After=local-fs.target network-online.target

[Service]
Type=notify
NotifyAccess=all
Restart=always
RestartSec=5

# There is a custom main process that kills all of the contained services.
KillMode=process
KillSignal=SIGTERM

ExecStartPre=/usr/bin/rm -f /var/lib/flatcar-oem-gce.img
ExecStartPre=/usr/bin/truncate -s 1G /var/lib/flatcar-oem-gce.img
ExecStartPre=/usr/sbin/mkfs.ext4 /var/lib/flatcar-oem-gce.img
ExecStartPre=/usr/bin/rm -rf /var/lib/flatcar-oem-gce
ExecStartPre=/usr/bin/mkdir -p /var/lib/flatcar-oem-gce
ExecStartPre=-/usr/bin/umount /var/lib/flatcar-oem-gce.img
ExecStartPre=/usr/bin/mount /var/lib/flatcar-oem-gce.img /var/lib/flatcar-oem-gce
ExecStartPre=/usr/bin/tar --directory=/var/lib/flatcar-oem-gce --extract --file=/usr/share/oem/flatcar-oem-gce.aci --strip-components=1 rootfs
ExecStartPre=/usr/bin/umount /var/lib/flatcar-oem-gce.img
ExecStart=/usr/bin/systemd-nspawn --keep-unit --register=no --link-journal=no \
    --machine=oem-gce --capability=CAP_NET_ADMIN --bind=/dev/log --bind=/run/systemd --tmpfs=/run/lock --bind=/etc --bind=/home --bind-ro=/usr/share/google-oslogin/nsswitch.conf \
    --read-only --volatile=overlay --image=/var/lib/flatcar-oem-gce.img /init.sh

[Install]
WantedBy=multi-user.target
